# packages/virtual/tvheadend42-standalone/system.d/tvheadend42.service
[Unit]
Description=TVHeadend42 Service
After=network-online.service

[Service]
ExecStart=-/bin/sh -c "exec sh /usr/share/tvheadend/bin/tvheadend42.start"
TimeoutStopSec=2
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
